# systemd/ai-worker-video.service
[Unit]
Description=AI worker video
After=network.target
[Service]
User=aiuser
WorkingDirectory=/srv/project
Environment=ROLE=video
ExecStart=/srv/project/bin/worker.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
