# systemd/ai-worker-compose.service
[Unit]
Description=AI worker compose
After=network.target
[Service]
User=aiuser
WorkingDirectory=/srv/project
Environment=ROLE=compose
ExecStart=/srv/project/bin/worker.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
