# systemd/ai-worker-audio.service
[Unit]
Description=AI worker audio
After=network.target
[Service]
User=aiuser
WorkingDirectory=/srv/project
Environment=ROLE=audio
ExecStart=/srv/project/bin/worker.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
